# debian/clamav-freshclam.prerm.in
#!/bin/sh

set -e

#COMMON-FUNCTIONS#

invoke-rc.d clamav-freshclam stop || true
#DEBHELPER#
exit 0
